# dist/ansible/roles/prometheus/templates/prometheus-reload.path
[Unit]
Description=Prometheus reload path

[Path]
PathChanged=/var/prometheus/config/prometheus.yml

[Install]
WantedBy=multi-user.target
