# dist/ansible/roles/prometheus/templates/prometheus-reload.service
[Unit]
Description=Prometheus reload service

[Service]
ExecStart=-/usr/bin/docker kill -s HUP prometheus
Type=oneshot
