# units/certificate-sync.service
[Unit]
Description=SSL Certificate Syncronization
# Dependencies
Requires=etcd.service

# Ordering
After=etcd.service

[Service]
ExecStartPre=etcdctl get /config/ssl > /etc/ssl/cert.crt
ExecStart=/bin/sh -c "while true; do etcdctl exec-watch  /config/ssl -- etcdctl get /config/ssl > /etc/ssl/cert.crt;done"

[X-Fleet]
Global=true
MachineOf=nginx.service
