# units/certificate-sync.service
[Unit]
Description=SSL Certificate Syncronization
# Dependencies
Requires=etcd.service

# Ordering
After=etcd.service

[Service]
ExecStart=etcdctl watch  /config/ssl
ExectStartPost=etcdctl get /config/ssl > /etc/ssl/cert.crt
Restart=always

[X-Fleet]
Global=true
MachineMetadata=frontend=true
