# units/nginx-reload.path
[Unit]
Description=NGINX reload path

[Path]
PathChanged=/etc/ssl
#TODO Add Local Config Path

[X-Fleet]
Global=true
MachineOf=nginx.service
