# units/nginx-reload.path
[Unit]
Description=NGINX reload path

[Path]
PathChanged=/etc/ssl
PathChanged=/var/nginx/nginx.conf

[X-Fleet]
Global=true
MachineOf=nginx.service
