# units/nginx-reload.path
[Unit]
Description=NGINX reload path

[Path]
PathChanged=/etc/ssl
PathChanged=/usr/var/nginx.conf

[X-Fleet]
Global=true
MachineMetadata=frontend=true
