# dist/ansible/roles/frontend/templates/acme-response-watcher.service
[Unit]
Description=Watches for distributed acme challenge responses
# Dependencies
Requires=etcd.service

# Ordering
After=etcd.service

[Service]
ExecStart=/usr/bin/etcdctl watch /acme/watched
ExecStartPost=-/usr/bin/systemctl start nginx-config-templater.service
Restart=always

[Install]
WantedBy=multi-user.target
