# dist/ansible/roles/frontend/templates/letsencrypt-renewal.timer
[Unit]
Description=Letsencrpyt renewal timer

[Timer]
OnCalendar=*-*-01 05:00:00
RandomizedDelaySec=60
