# dist/ansible/roles/frontend/templates/nginx-reload.path
[Unit]
Description=NGINX reload path

[Path]
PathChanged=/var/nginx/nginx.conf
PathChanged=/var/ssl/fullchain.pem

[Install]
WantedBy=multi-user.target
