# dist/ansible/roles/frontend/templates/nginx-reload.service
[Unit]
Description=NGINX reload service

[Service]
ExecStart=-/usr/bin/docker kill -s HUP nginx
Type=oneshot
