# dist/ansible/roles/frontend/templates/route-discovery-watcher.service
[Unit]
Description=Watches for nginx routes
# Dependencies
Requires=etcd.service

# Ordering
After=etcd.service

# Restart when dependency restarts
PartOf=etcd.service

[Service]
ExecStartPre=-/bin/sh -c '/usr/bin/etcdctl mk /route_discovery/watched "$(date +%s%N)"'
ExecStart=/usr/bin/etcdctl watch /route_discovery/watched 
ExecStartPost=-/usr/bin/systemctl start nginx-config-templater.service
Restart=always

[Install]
WantedBy=multi-user.target
